# configs/systemd/tiny-dfr.service.d/override.conf
# Managed by dotfiles — installed by script/touchbar/arch/setup.sh
#
# Intel T2 Mac fixes for the upstream tiny-dfr.service:
#
# 1. BindsTo= : upstream binds to /dev/tiny_dfr_display and
#    /dev/tiny_dfr_backlight, which are Apple-Silicon device nodes. On Intel T2
#    the Touch Bar display is a USB appletbdrm device, so those .device units
#    never activate and the daemon's start job fails with a 'dependency' error.
#    Clearing BindsTo lets tiny-dfr start and open whatever display is present.
#
# 2. graphical.target : the daemon must NOT be pulled into multi-user.target
#    while ordered After=graphical.target — graphical.target already starts
#    After=multi-user.target, so that combination is an ordering cycle and
#    systemd silently drops the unit at boot (blank Touch Bar). Order after and
#    install into graphical.target instead.

[Unit]
# Drop the Apple-Silicon device bindings (display + backlight).
BindsTo=
After=graphical.target

[Install]
WantedBy=graphical.target
