# configs/rift/com.nest.appearance-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nest.appearance-watcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>$HOME/dev/dotfiles/configs/rift/appearance-watcher.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>~/Library/Preferences/.GlobalPreferences.plist</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/appearance-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/appearance-watcher.log</string>
</dict>
</plist>
